# yak-upstart.conf
#!upstart
description "node.js server"
author      "Someone that likes Node"
 
# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on (local-filesystems and net-device-up)
stop on shutdown
 
# Automatically Respawn:
respawn
respawn limit 99 5
 
script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/home/ubuntu"
    export NODE_ENV="production" 
    exec /usr/bin/node /home/ubuntu/sites/yak/site/bin/www >> /var/log/yak-upstart.log 2>&1
end script
 
post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script

